# net-misc/kea/files/kea-initd
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="kea dhcp services"

dhcp4_command="/usr/sbin/kea-dhcp4"
dhcp6_command="/usr/sbin/kea-dhcp6"
ddns_command="/usr/sbin/kea-dhcp-ddns"
control_command="/usr/sbin/kea-ctrl-agent"
config_file="/etc/kea/kea.conf"
dhcp4_pidfile="/var/lib/kea/kea.kea-dhcp4.pid"
dhcp6_pidfile="/var/lib/kea/kea.kea-dhcp6.pid"
ddns_pidfile="/var/lib/kea/kea.kea-dhcp-ddns.pid"
control_pidfile="/var/lib/kea/kea.kea-ctrl-agent.pid"

depend() {
	use net mysql postgresql
}

start_pre() {
	if [ ! -f "${config_file}" ] ; then
		eerror "Please create a ${config_file} config file."
		return 1
	fi

	if ${DHCP4:-false} ; then
		if ! ${dhcp4_command} -t ${config_file} 1>/dev/null 2>/dev/null ; then
			eerror "Error in config file."
			return 1
		fi
	fi
	if ${DHCP6:-false} ; then
		if ! ${dhcp6_command} -t ${config_file} 1>/dev/null 2>/dev/null ; then
			eerror "Error in config file."
			return 1
		fi
	fi
	if ${CONTROL:-false} ; then
		if ! ${control_command} -t ${config_file} 1>/dev/null 2>/dev/null ; then
			eerror "Error in config file."
			return 1
		fi
	fi
}

start() {
	if ${DHCP4:-false} ; then
		ebegin "Starting kea dhcp4 service"
		start-stop-daemon -b -p ${dhcp4_pidfile} \
			-x ${dhcp4_command} -- -c ${config_file} \
			|| return 1
		eend $?
	fi
	if ${DHCP6:-false} ; then
		ebegin "Starting kea dhcp6 service"
		start-stop-daemon -b -p ${dhcp6_pidfile} \
			-x ${dhcp6_command} -- -c ${config_file} \
			|| return 1
		eend $?
	fi
	if ${DDNS:-false} ; then
		ebegin "Starting kea ddns service"
		start-stop-daemon -b -p ${ddns_pidfile} \
			-x ${ddns_command} -- -c ${config_file} \
			|| return 1
		eend $?
	fi
	if ${CONTROL:-false} ; then
		ebegin "Starting kea control service"
		start-stop-daemon -b -p ${control_pidfile} \
			-x ${control_command} -- -c ${config_file} \
			|| return 1
		eend $?
	fi
}

stop() {
	if ${DHCP4:-false} ; then
		ebegin "Stopping kea dhcp4 service"
		start-stop-daemon --stop -p ${dhcp4_pidfile} \
			|| return 1
		eend $?
	fi
	if ${DHCP6:-false} ; then
		ebegin "Stopping kea dhcp6 service"
		start-stop-daemon --stop -p ${dhcp6_pidfile} \
			|| return 1
		eend $?
	fi
	if ${DDNS:-false} ; then
		ebegin "Stopping kea ddns service"
		start-stop-daemon --stop -p ${ddns_pidfile} \
			|| return 1
		eend $?
	fi
	if ${CONTROL:-false} ; then
		ebegin "Stopping kea control service"
		start-stop-daemon --stop -p ${control_pidfile} \
			|| return 1
		eend $?
	fi
}
