# sys-cluster/seaweedfs/files/weed-filer.initd
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="SeaweedFS filer service"
pidfile=${pidfile:-"/run/${RC_SVCNAME}.pid"}
command="/usr/bin/weed"
command_user="${command_user:-"seaweedfs:seaweedfs"}"
command_args="${command_args:-"filer -localSocket=/run/seaweedfs/${RC_SVCNAME}.sock"}"
command_background="true"
logfile="${logfile:-"/var/log/seaweedfs/${RC_SVCNAME}.log"}"
start_stop_daemon_args="--user ${command_user%:*} --group ${command_user#*:} --stdout ${logfile} --stderr ${logfile}"

: "${supervisor:=supervise-daemon}"
: "${respawn_delay:=5}"
: "${respawn_max:=10}"
: "${respawn_period:=60}"

start_pre() {
	checkpath -d -o ${command_user%:*} -q /run/seaweedfs
}

depend() {
	need net
	after weed-master
}
